# loadkeys.service
[Unit]
Description=loadkeys

[Service]
ExecStart=/usr/local/bin/loadkeys-unit

[Install]
WantedBy=multi-user.target