# zram.service
[Unit]
Description=Use zram to get more ram-like swap space.

[Service]
ExecStart=/home/xyh/conf/zram

[Install]
WantedBy=multi-user.target